# .config/systemd/user/mbsync.service
[Unit]
Description=mbsync

[Service]
Type=oneshot
Environment=GNUPGHOME=/home/b/sync/gnupghome
ExecStart=/usr/bin/mbsync -Va
